# void-pkgs/etc/sv/runsvdir-wc/run
#!/bin/sh

export USER="wc"
export HOME="/home/wc"

groups="$(id -Gn "$USER" | tr ' ' ':')"
svdir="$HOME/.local/service"

exec chpst -u "$USER:$groups" runsvdir "$svdir"
