# void-cfg/.local/service/pulseaudio/run
#!/bin/sh
exec pulseaudio --daemonize=no --exit-idle-time=-1
