# void-cfg/.local/service/gpg-agent/run
#!/bin/sh

exec /usr/bin/gpg-agent --supervised
